# ansible/roles/sentry/templates/sentry.supervisord.conf
[program:sentry-web]
directory="{{ app_home }}"
command="{{ app_env }}/bin/sentry start http"
autostart=true
autorestart=true
redirect_stderr=true
